# build/proxy-gateway/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:tailscaled]
command=/usr/local/bin/tailscaled --tun=userspace-networking --state=/var/lib/tailscale/tailscaled.state --socket=/var/run/tailscale/tailscaled.sock
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/tailscaled.err.log
stdout_logfile=/var/log/tailscaled.out.log
environment=PATH="/usr/local/bin:/usr/bin:/bin"

[program:tailscale-auth]
command=/usr/local/bin/setup-tailscale-proxy.sh
user=root
autostart=true
autorestart=false
startretries=5
stderr_logfile=/var/log/tailscale-auth.err.log
stdout_logfile=/var/log/tailscale-auth.out.log
depends_on=tailscaled

[program:mitmproxy]
command=/usr/local/bin/start-mitmproxy.sh
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/mitmproxy.err.log
stdout_logfile=/var/log/mitmproxy.out.log
depends_on=tailscale-auth

[program:proxy-routing]
command=/usr/local/bin/setup-proxy-routing.sh
user=root
autostart=true
autorestart=true
stderr_logfile=/var/log/proxy-routing.err.log
stdout_logfile=/var/log/proxy-routing.out.log
depends_on=mitmproxy
